# docs/setup/systemd/vulnhunter.service.example
[Unit]
Description=VulnHunter governed web application
After=network.target
ConditionPathExists=/srv/vulnhunter/current/manage.py

[Service]
Type=notify
User=vulnhunter
Group=vulnhunter
WorkingDirectory=/srv/vulnhunter/current
EnvironmentFile=/etc/vulnhunter/vulnhunter.env
RuntimeDirectory=vulnhunter
RuntimeDirectoryMode=0750
UMask=0077
ExecStartPre=/srv/vulnhunter/venv/bin/python manage.py check
ExecStartPre=/srv/vulnhunter/venv/bin/python manage.py migrate --check
ExecStart=/srv/vulnhunter/venv/bin/gunicorn --config config/deployment/gunicorn.conf.py vulnhunter.web.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGTERM
TimeoutStopSec=45
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/vulnhunter/state /srv/vulnhunter/evidence /srv/vulnhunter/log /srv/vulnhunter/static

[Install]
WantedBy=multi-user.target
